# arch-sense.service
[Unit]
Description=Arch-Sense — Apply saved keyboard RGB settings
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/arch-sense --apply
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
